# common/.config/systemd/user/audio-share.service
[Unit]
Description=Audio share server

[Service]
ExecStart=/usr/bin/as-cmd -b 0.0.0.0

[Install]
WantedBy=default.target
