# common/.config/systemd/user/audio-share.service
[Unit]
Description=Audio Share Server

[Service]
ExecStart=/usr/bin/as-cmd -b 0.0.0.0 -e 56

[Install]
WantedBy=default.target
